# udev/etc/dev/rules.d/97-usb-camera-trigger.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idProduct}!="", RUN+="/usr/local/bin/take-picture.sh"
