# udev/etc/dev/rules.d/98-power-notify.rules
SUBSYSTEM=="power_supply", KERNEL=="AC0", ATTR{online}=="1", RUN+="/usr/local/bin/ac-notify.sh"
