# udev/etc/dev/rules.d/99-mouse-notify.rules
ACTION=="add", SUBSYSTEM=="input", ATTR{name}=="Primax Kensington Eagle Trackball", ENV{ID_INPUT_MOUSE}=="1", RUN+="/usr/local/bin/mouse-notify.sh"
